# TerminalApp/SoyehtMac/Library/LaunchAgents/com.soyeht.engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!-- LaunchAgent registered by SMAppService.agent(plistName:). -->
<dict>
	<key>Label</key>
	<string>com.soyeht.engine</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>SOYEHT_DIR="$HOME/Library/Application Support/Soyeht"; ENGINE_DIR="$SOYEHT_DIR/engine"; export ADMIN_PORT="8892"; export ADDR="0.0.0.0:8892"; export THEYOS_BIN_DIR="$ENGINE_DIR"; export THEYOS_SQLITE_DB="$SOYEHT_DIR/theyos.db"; export THEYOS_RATELIMIT_DB="$SOYEHT_DIR/ratelimit.db"; export THEYOS_CONVERSATIONS_DIR="$SOYEHT_DIR/conversations"; export THEYOS_SKIP_LEGACY_MIGRATION="1"; export THEYOS_FORCE_SOFTWARE_KEYS="1"; export THEYOS_VMRUNNER_RS_BIN="$ENGINE_DIR/vmrunner_macos_ipc"; export THEYOS_STORE_RS_BIN="$ENGINE_DIR/store-ipc"; export THEYOS_TERMINAL_RS_BIN="$ENGINE_DIR/terminal-ipc"; export THEYOS_SSH_CTL="$ENGINE_DIR/theyos-ssh"; export THEYOS_APNS_KEY_PATH="$SOYEHT_DIR/apns.p8"; export THEYOS_APNS_KEY_ID="5FPYV735V4"; export THEYOS_APNS_TEAM_ID="W7677A5BK2"; export THEYOS_APNS_TOPIC="com.soyeht.app"; exec "$ENGINE_DIR/theyos-engine"</string>
	</array>

	<key>EnvironmentVariables</key>
	<dict>
		<key>THEYOS_APNS_KEY_ID</key>
		<string>5FPYV735V4</string>
		<key>THEYOS_APNS_TEAM_ID</key>
		<string>W7677A5BK2</string>
		<key>THEYOS_APNS_TOPIC</key>
		<string>com.soyeht.app</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardErrorPath</key>
	<string>/tmp/soyeht-engine.log</string>

	<key>StandardOutPath</key>
	<string>/tmp/soyeht-engine.log</string>
</dict>
</plist>
